# data/ippd.socket
[Unit]
Description=IPP Daemon socket

[Socket]
ListenStream=631

[Install]
WantedBy=sockets.target
