# data/ippd.service.in
[Unit]
Description=IPP Daemon

[Service]
ExecStart=@libexecdir@/ippd
StandardOutput=syslog
StandardError=syslog
