# data/printerd.service.in
[Unit]
Description=Printer Daemon

[Service]
Type=dbus
BusName=org.freedesktop.printerd
ExecStart=@privlibdir@/printerd
